# ops/templates/aragonteam-watchdog.timer
[Unit]
Description=aragonteam Watchdog Timer

[Timer]
OnBootSec=120
OnUnitActiveSec=60s
AccuracySec=5s

[Install]
WantedBy=timers.target
